# dev_setup/cookbooks/nfs_server/templates/default/nfs-umount.conf
description "Ensures all nfs shares in fstab are umounted"

start on runlevel [06]

emits nfsshares-mounted
emits nfsshares-not-mounted

script
   exec /sbin/ensure-nfs-umounted
end script

