# dev_setup/cookbooks/deployment/templates/default/etc_issue.conf.erb
# etc issues volatile messages
#
# This service maintains a directory of temporary messages displayed
# on the console's issue from the point the system is
# started until it is shut down again.

description "etc issues volatile messages"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

expect fork

pre-start script
    # put the folder that will contain the messages in place
    [ -d /etc/issue.d/volatile.d ] && rm -rf /etc/issue.d/volatile.d
    mkdir -p /etc/issue.d/volatile.d
    touch /etc/issue.d/volatile.d/K00_vcap_welcome
    echo "<%= node[:deployment][:etc_issue_msg] %>" > /etc/issue.d/volatile.d/K00_vcap_welcome
    chmod -R go+rw /etc/issue.d/volatile.d
end script

post-stop script
    # remove the folders that contain the messages
    [ -d /etc/issue.d/volatile.d ] && rm -rf /etc/issue.d/volatile.d
end script

script
    /etc/issue_update &
end script
