# dev_setup/cookbooks/deployment/templates/default/vcap_shutdown.conf.erb
# vcap_shutdown - task to shutdown vcap and clean-up the pid files.
#
description     "vcap shutdown"

start on runlevel [06]

task

script
 set +e
 LOG_FILE=/var/log/vcap_reconfig.log
 touch $LOG_FILE
 chmod o+r $LOG_FILE
 monit quit
 sudo -i -u <%= node[:deployment][:user] %> <%= node[:deployment][:vcap_exec] %> stop

 echo "post-stop for this vcap_reconfig is called to clean-up the pid files" >> $LOG_FILE
 [ -d "/var/vcap/sys/run" ] && rm -rf /var/vcap/sys/run/*

end script

