# dev_setup/cookbooks/nats_server/templates/default/nats_server.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          nats_server
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       Starts the nats message bus server
### END INIT INFO

# Author: VMware
PATH=<%= File.join(node[:ruby][:path], "bin") %>:/sbin:/usr/sbin:/bin:/usr/bin

DESC="nats-server"
NAME=nats-server
DAEMON_NAME=nats-server
DAEMON=<%= File.join(node[:ruby][:path], "bin", "nats-server") %>
DAEMON_ARGS="-d -c <%= node[:nats_server][:config] %>"

PIDFILE=/var/run/nats_server.pid
SCRIPTNAME=/etc/init.d/$NAME

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

# Read configuration variable file if it is present
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.0-6) to ensure that this file is present.
. /lib/lsb/init-functions

#
# Function that starts the daemon/service
#
do_start()
{
    # Return
    #   0 if daemon has been started
    #   1 if daemon was already running
    #   2 if daemon could not be started
    status_of_proc "$DAEMON" "$NAME" > /dev/null
    if [ "$?" = 0 ]; then
      echo "The nats-server is already running"
      exit 1
    fi
    start-stop-daemon --start --quiet --pidfile $PIDFILE --name $DAEMON_NAME --exec $DAEMON -- \
        $DAEMON_ARGS || return 2
    #LOG_FILE=/var/log/nats_server_ip_changes.log
    #echo "Started nats from /etc/init.d; the content of $PIDFILE "`cat $PIDFILE` >> $LOG_FILE
    #ppid=`ps -p $$ -o ppid=`
    #echo "ppid=$ppid" >> $LOG_FILE
}

#
# Function that stops the daemon/service
#
do_stop()
{
    # Return
    #   0 if daemon has been stopped
    #   1 if daemon was already stopped
    #   2 if daemon could not be stopped
    #   other if a failure occurred
    status_of_proc "$DAEMON" "$NAME" > /dev/null
    if [ "$?" != 0 ]; then
      echo "The nats-server is not running."
      if [ -f $PIDFILE ]; then
        echo "Warning removing the orphan pid file $PIDFILE"
        rm -f $PIDFILE
      fi
      exit 1
    fi

    start-stop-daemon --stop --retry=TERM/30/KILL/5 --pidfile $PIDFILE
    RETVAL="$?"
    [ "$RETVAL" = 2 ] && "Warn: unable to kill the process based on the pid. Trying with name."#return 2
    [ "$RETVAL" = 0 ] && oknodo="--quiet --oknodo"
    start-stop-daemon --stop $oknodo --oknodo --retry=0/30/KILL/5 --name $DAEMON_NAME
    [ "$?" = 2 ] && return 2
    rm -f $PIDFILE
    return "$RETVAL"
}

#
# Function that sends a SIGHUP to the daemon/service
#
do_reload() {
    #
    # If the daemon can reload its configuration without
    # restarting (for example, when it is sent a SIGHUP),
    # then implement that here.
    #
    start-stop-daemon --stop --signal 1 --quiet --pidfile $PIDFILE --name $DAEMON_NAME
    return 0
}

case "$1" in
  start)
    [ "$VERBOSE" != no ] && log_daemon_msg "Starting $DESC" "$NAME"
    do_start
    case "$?" in
        0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
        2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
    esac
    ;;
  stop)
    [ "$VERBOSE" != no ] && log_daemon_msg "Stopping $DESC" "$NAME"
    do_stop
    case "$?" in
        0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
        2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
    esac
    ;;
  status)
       status_of_proc "$DAEMON" "$NAME" && exit 0 || exit $?
       ;;
  restart|force-reload)
    log_daemon_msg "Restarting $DESC" "$NAME"
    do_stop
    case "$?" in
      0|1)
        do_start
        case "$?" in
            0) log_end_msg 0 ;;
            1) log_end_msg 1 ;; # Old process is still running
            *) log_end_msg 1 ;; # Failed to start
        esac
        ;;
      *)
          # Failed to stop
        log_end_msg 1
        ;;
    esac
    ;;
  *)
    echo "Usage: $SCRIPTNAME {start|stop|status|restart|force-reload}" >&2
    exit 3
    ;;
esac

:
